# storybrain.service
[Unit]
Description=Gunicorn instance to serve StoryBrain AI
After=network.target

[Service]
# Replace 'ubuntu' with your VPS username
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/tool
Environment="PATH=/home/ubuntu/tool/venv/bin"
# Adjust workers based on RAM (1GB = 2, 2GB = 4)
ExecStart=/home/ubuntu/tool/venv/bin/gunicorn -c gunicorn_conf.py main:app

[Install]
WantedBy=multi-user.target
