# storybrain-ai.service
[Unit]
Description=StoryBrain AI — Multi-Tool Website
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/storybrain-ai

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=HOME=/root

ExecStart=/opt/storybrain-ai/venv/bin/gunicorn app.main:app -c /opt/storybrain-ai/gunicorn_conf.py

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
