# storybrain-ai.service
[Unit]
Description=StoryBrain AI — Multi-Tool Website
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory={{APP_DIR}}

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=HOME=/root
Environment=U2NET_HOME={{APP_DIR}}/.u2net

ExecStart={{APP_DIR}}/venv/bin/gunicorn app.main:app -c {{APP_DIR}}/gunicorn_conf.py

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
